# docker/resources/cron/application-extract.crontab
# Every quarter-past-the-hour during the development phase
15	*	*	*	*	cd /app && node src/processor.js