# docker/resources/cron/application-extract.crontab
*	*	*	*	*	cd /app && node src/processor.js