# docker/resources/cron/application-extract.crontab
*	*	*	*	*	cd /usr/src/app && node src/processor.js