# docker/resources/cron/repdata-extract.crontab
# Every quarter-to-the-hour during the development phase
45	*	*	*	*	cd /app && node src/processor.js